# infra/templates/systemd/load-reth@.service.j2
[Unit]
Description=load-reth (%i)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=2
LimitNOFILE=1048576

WorkingDirectory=/opt/loadnet/current

EnvironmentFile=/opt/loadnet/current/bundle/private/env/load-reth-%i.env

ExecStartPre=/bin/bash -lc 'test -f "/opt/loadnet/current/bundle/public/genesis.json"'
ExecStartPre=/bin/bash -lc 'test -f "/opt/loadnet/current/$LOAD_RETH_P2P_KEY_PATH"'
ExecStartPre=/bin/bash -lc 'rm -f "/run/load-reth/%i/engine.ipc"'
ExecStartPre=/bin/bash -lc 'docker rm -f load-reth-%i >/dev/null 2>&1 || true'

ExecStart=/bin/bash -lc '/usr/bin/docker run --rm --name load-reth-%i --network host --platform=linux/amd64 \
  --cap-drop=ALL \
  --security-opt=no-new-privileges \
  --env HOME=/data/load-reth \
  --env XDG_CACHE_HOME=/data/load-reth/.cache \
  --log-driver={{ loadnet_docker_log_driver | default("local") }} \
  --log-opt max-size={{ loadnet_docker_log_max_size | default("50m") }} \
  --log-opt max-file={{ loadnet_docker_log_max_file | default("5") }} \
  --user $LOAD_RETH_UID:$LOAD_RETH_GID \
  -v /var/lib/load-reth/%i:/data/load-reth \
  -v /opt/loadnet/current/bundle/public:/assets:ro \
  -v /opt/loadnet/current/$LOAD_RETH_P2P_KEY_PATH:/run/load-reth/%i/p2p.key:ro \
  -v /run/load-reth/%i:/run/load-reth/%i \
  $LOAD_RETH_IMAGE node \
    --datadir=/data/load-reth \
    --chain=/assets/genesis.json \
    --http --http.addr={{ loadnet_el_http_bind | default("0.0.0.0") }} --http.port=$LOAD_RETH_HTTP_PORT --http.api={{ loadnet_el_http_api | default("eth,net,web3,txpool") }} \
    --authrpc.addr=127.0.0.1 --authrpc.port=$LOAD_RETH_AUTHRPC_PORT \
    --metrics=127.0.0.1:$LOAD_RETH_METRICS_PORT \
    --port=$LOAD_RETH_P2P_PORT --discovery.port=$LOAD_RETH_P2P_PORT \
    --nat=extip:$LOAD_RETH_PUBLIC_IP \
    --p2p-secret-key=/run/load-reth/%i/p2p.key \
    --bootnodes=$LOAD_RETH_BOOTNODES \
    --auth-ipc --auth-ipc.path=$LOAD_RETH_ENGINE_IPC_PATH'

ExecStop=/usr/bin/docker stop -t 10 load-reth-%i
ExecStartPost=/bin/bash -lc 'for i in {1..60}; do if [ -S "/run/load-reth/%i/engine.ipc" ]; then chgrp {{ loadnet_ultramarine_gid | default(10002) }} "/run/load-reth/%i/engine.ipc"; chmod 660 "/run/load-reth/%i/engine.ipc"; exit 0; fi; sleep 1; done; echo "engine IPC not ready for permission fix: /run/load-reth/%i/engine.ipc" >&2; exit 0'

[Install]
WantedBy=multi-user.target
